# infra/gunicorn.service
[Unit]
Description=IRG GDP Gunicorn application server
After=network.target

[Service]
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/irg_gdp/backend
EnvironmentFile=/home/ec2-user/irg_gdp/backend/.env
ExecStart=/home/ec2-user/.local/bin/gunicorn \
    wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 3 \
    --timeout 60 \
    --access-logfile /home/ec2-user/gunicorn-access.log \
    --error-logfile /home/ec2-user/gunicorn-error.log
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
